# StartupServices/plexRefresh/plex_refresh-template.service
[Unit]
Description=Refresh Plex metadata
Wants=network-online.target
After=network-online.target mount_drives.service plex_rename_folders.service update_plex_permissions.service
Requires=mount_drives.service plex_rename_folders.service update_plex_permissions.service

# Ensure all media directories are mounted before running
RequiresMountsFor=/mnt/media/Media
RequiresMountsFor=/mnt/media/Media/Backup
RequiresMountsFor=/mnt/media/Media/Movies
RequiresMountsFor=/mnt/media/Media/Music
RequiresMountsFor=/mnt/media/Media/Other\x20Videos
RequiresMountsFor=/mnt/media/Media/Photos
RequiresMountsFor=/mnt/media/Media/TV\x20Series

# Only run if the main mount point exists
ConditionPathIsMountPoint=/mnt/media/Media

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /usr/local/bin/plex_refresh.py
StandardOutput=append:/var/log/plex_refresh.log
StandardError=append:/var/log/plex_refresh.log
RemainAfterExit=true
User=root
Group=root
Environment="PYTHONUNBUFFERED=1" "LANG=C.UTF-8" "LC_ALL=C.UTF-8"

[Install]
WantedBy=multi-user.target