# ShutdownServices/rsyncBackup/rsync_backup-template.service
[Unit]
Description=rsync_backup
DefaultDependencies=no
RequiresMountsFor=/mnt/backup/Backups
After=local-fs.target mount_drives.service
# If you don't have mount_drives.service, remove it from the After= line above.
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecStop=/usr/bin/python3 -u /usr/local/bin/rsync_backup.py
User=root
Group=root
Environment="PYTHONUNBUFFERED=1" "LANG=C.UTF-8" "LC_ALL=C.UTF-8"
StandardOutput=append:/var/log/rsync_backup.log
StandardError=append:/var/log/rsync_backup.log
TimeoutStopSec=30min

[Install]
WantedBy=multi-user.target